# deploy/oracle/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw gateway (SyncScript) on loopback
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/syncscript
Environment=NODE_ENV=production
# Requires: Node 22+ and `openclaw` on PATH (see deploy/oracle/vm/install-stack-services.sh)
ExecStart=/usr/bin/openclaw gateway run --port 18789 --bind loopback
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
